# systemd/asac.service
[Unit]
Description=ASAC - Sistema de Gestao de Materiais (FastAPI/Uvicorn)
After=network.target

[Service]
Type=simple
User=fabricio
Group=fabricio
WorkingDirectory=/home/fabricio/Projetos/Python/yolo-barcode
ExecStart=/home/fabricio/Projetos/Python/yolo-barcode/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --ssl-keyfile certs/key.pem --ssl-certfile certs/cert.pem
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
